# new-api.service
# File path: /etc/systemd/system/ai-bridge.service
# sudo systemctl daemon-reload
# sudo systemctl start ai-bridge
# sudo systemctl enable ai-bridge
# sudo systemctl status ai-bridge
[Unit]
Description=One API Service
After=network.target

[Service]
User=ubuntu  # 注意修改用户名
WorkingDirectory=/path/to/ai-bridge  # 注意修改路径
ExecStart=/path/to/ai-bridge/ai-bridge --port 3000 --log-dir /path/to/ai-bridge/logs  # 注意修改路径和端口号
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
